# etc/conf.d/django.conf
[program:gopherairtime-django]
command = /srv/wcl/prod/gopherairtime/ve/bin/gunicorn --bind 0.0.0.0:8000 gopherairtime.wsgi
directory = /srv/wcl/prod/gopherairtime
environment = PYTHONPATH='.'
stdout_logfile = ./logs/%(program_name)s_%(process_num)s.log
stderr_logfile = ./logs/%(program_name)s_%(process_num)s.log
