# etc/conf.d/rabbitmq-server.conf
[program:rabbitmq-server]
command=/usr/sbin/rabbitmq-server
stdout_events_enabled=true
stderr_events_enabled=true
